# meta-nxp-matter-advanced/recipes-core/udev/udev-rules-imx/98-btnxpuart.rules
ACTION=="change", ENV{BTNXPUART_DEV}=="*", ENV{BTNXPUART_STATE}=="*", RUN+="/etc/udev/scripts/btnxpuart.sh $env{BTNXPUART_DEV} $env{BTNXPUART_STATE}"
